# etc/systemd/system/expandrootfs.service
[Unit]
Description=Expand rootfs

[Service]
ExecStart=/etc/solydxk/rpi/expandrootfs

[Install]
WantedBy=multi-user.target
